# recovery/thinxai-telegram.service
[Unit]
Description=ThinxAI Telegram Bridge
After=network.target

[Service]
Type=simple
User=jdlongmire
Group=jdlongmire
WorkingDirectory=/media/jdlongmire/Macro-Drive-2TB/GitHub_Repos/thinx/thinxai-telegram
Environment=PATH=/home/jdlongmire/.nvm/versions/node/v24.13.0/bin:/usr/local/bin:/usr/bin:/bin:/home/jdlongmire/.local/bin
ExecStart=/usr/bin/python3 bridge.py
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thinxai-telegram

[Install]
WantedBy=multi-user.target
