# recovery/thinxai-telegram.service
[Unit]
Description=ThinxAI Telegram Bridge
After=network.target

[Service]
Type=simple
User=$USER
Group=$USER
# Update WorkingDirectory to your thinxai-telegram location:
WorkingDirectory=/path/to/thinx/thinxai-telegram
# Update PATH to include your nvm node and local bin:
Environment=PATH=/home/$USER/.nvm/versions/node/v24.x.x/bin:/usr/local/bin:/usr/bin:/bin:/home/$USER/.local/bin
ExecStart=/usr/bin/python3 bridge.py
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thinxai-telegram

[Install]
WantedBy=multi-user.target
